# openshift-node-checker.service
[Unit]
Description=OpenShift Node Checker

[Service]
Type=simple
ExecStart=/usr/bin/python3 /path/to/script.py
Restart=on-failure

[Install]
WantedBy=multi-user.target
